# hostname/templates/default/rc.local.erb
#!/bin/sh -e

bash /root/set-hostname.sh
exit 0
